# wsl/profile.d/bad_colors.bash
# These are the colors based on how they show up in my terminal

export reset=$'\e[0m'
export bold=$'\e[1m'
export dim=$'\e[2m'
export invert=$'\e[7m'
export invis=$'\e[8m'
export black=$'\e[30m'
export red=$'\e[31m'
export green=$'\e[32m'
export purple=$'\e[34m'
export pink=$'\e[35m'
export indigo=$'\e[36m'
export grey=$'\e[90m'
export yellow=$'\e[93m'
export cyan=$'\e[96m'

export blue=$'\e[38;5;018m'
export teal=$'\e[38;5;122m'
export cyan=$'\e[38;5;074m'
export orange=$'\e[38;5;173m'

# tends to show up more orangey
# export yellow=$'\e[33m'
