# wsl/profile.d/paths.sh
#!/bin/bash

export C
export dotfiles="$ambient/wsl"
PATH="$PATH:$dotfiles/bin:/trash/bin:/trash/been/bin:$C/Windows:$C/Windows/system32"
# Get some windows paths going on
