# wsl/profile.d/paths.sh
export C
PATH="/home/linuxbrew/.linuxbrew/bin:$PATH:$dotfiles/bin"
export JAVA_HOME="$C/Program Files/OpenJDK/jdk-11.0.2/"
# Get some windows paths going on
PATH+=:$C/Windows:$C/Windows/system32:$JAVA_HOME/bin
