# wsl/profile.d/paths.bash
export PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:$dotfiles/bin"
export C=/mnt/c
