# wsl/profile.d/paths.bash
export PATH="/home/linuxbrew/.linuxbrew/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:$dotfiles/bin"
export C=/mnt/c
# Get some windows paths going on
PATH+=:$C/Windows:$system32:$C/Go/bin
