# wsl/profile.d/paths.bash
source $dotfiles/environment
PATH="/home/linuxbrew/.linuxbrew/bin:$PATH:$dotfiles/bin"
# Get some windows paths going on
PATH+=:$C/Windows:$C/Windows/system32
